# init-system-helpers/wheezy/v8-devel/debian/init-system-helpers.install
script/deb-systemd-helper usr/bin
script/deb-systemd-invoke usr/bin
